# scripts/com.naly.cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.naly.cron</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/thomasjeon/GitHub/@dev.realbits/Naly/scripts/local-cron.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/thomasjeon/GitHub/@dev.realbits/Naly/logs/launchd-out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/thomasjeon/GitHub/@dev.realbits/Naly/logs/launchd-err.log</string>
</dict>
</plist>
